# blobcatd.service
[Unit]
Description=a systemd service for blobcatd

[Service]
ExecStart=blobcatd

[Install]
WantedBy=default.target
